# spark/supervisord.conf

[program:spark]
user=hduser
command=/opt/spark/sbin/start-all.sh